# cpp/data/installer/boot_notification.service
[Unit]
Description=Waity Service
After=systend-user-sessions.service

[Service]
Type=simple
ExecStart=/usr/local/bin/boot_notification.sh
